# scripts/systemd/authorclaw-build.timer
[Unit]
Description=Poll for the BookClaw build_now sentinel once a minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=15s

[Install]
WantedBy=timers.target
